# services/hopper-driver.service
[Unit]
Description=Hopper leg motor driver (AK60 via CAN)
After=canable.service
Requires=canable.service
# Wheel bus (can1) is OPTIONAL: order after it so DmWheelController::init sees
# can1 when the adapter is present, but never block hop-only operation on it.
After=canable2.service
Wants=canable2.service

[Service]
Type=simple
WorkingDirectory=/home/nvidia/Hopper_srbRL
ExecStart=/home/nvidia/Hopper_srbRL/build/hopper_driver
# Runs as root (CAN access / realtime), same as the old `sudo ./build/hopper_driver`.
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
