# services/px4-bridge.service
[Unit]
Description=PX4 USB MAVLink bridge (RETIRED: props now go over TELEM2/DDS via px4-dds-bridge)
# 2026-07 "Plan B": prop commands moved to the TELEM2/DDS link (px4_dds_bridge.py
# publishes VehicleCommand DO_SET_ACTUATOR on /fmu/in/vehicle_command). This USB
# service is kept ONLY for bench debugging; connect_and_start.sh stops+disables it.
# If you ever re-enable it, run with --no-props too (or stop px4-dds-bridge),
# otherwise TWO prop streams fight over the actuators.
After=network-online.target canable.service
Wants=network-online.target

[Service]
Type=simple
User=nvidia
Group=dialout
WorkingDirectory=/home/nvidia/Hopper_srbRL/pixhawk
Environment=HOME=/home/nvidia
ExecStart=/bin/bash /home/nvidia/Hopper_srbRL/pixhawk/run_px4_bridge.sh --no-imu --no-props
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
